# Installer/build/it.ipzs.CIE-ID-Bar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>it.ipzs.CIE-ID-Bar</string>
    <key>ProgramArguments</key>
    <array>
	<string>/usr/bin/open</string>
	<string>-n</string>
      	<string>/Applications/CIE ID Bar.app</string>
    </array>
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
        <string>/dev/null</string>
    <key>StandardOutPath</key>
        <string>/dev/null</string>
</dict>
</plist>